# usr/lib/systemd/system/dhcp6r.service
[Unit]
Description=Wide DHCPv6 relay agent
Documentation=man:dhcp6relay(8)

[Service]
EnvironmentFile=-/etc/sysconfig/dhcp6r
ExecStart=/usr/sbin/dhcp6relay -f $DHCP6RARGS

[Install]
WantedBy=network.target

